# xow.service
[Unit]
Description=Xbox One Wireless Dongle Driver

[Service]
ExecStart=/usr/local/bin/xow

[Install]
WantedBy=multi-user.target
